# deployment/configs/china-car-parts-api.service
[Unit]
Description=China Car Parts API - Production Environment
Documentation=https://github.com/your-org/china-car-parts
After=network.target postgresql.service
Wants=postgresql.service
Requires=network.target

[Service]
Type=exec
User=partsbot
Group=partsbot
WorkingDirectory=/opt/china-car-parts
Environment=PATH=/opt/china-car-parts/venv/bin
Environment=APP_ENV=production
EnvironmentFile=/opt/china-car-parts/.env
ExecStart=/opt/china-car-parts/venv/bin/uvicorn app.api.main:app --host 0.0.0.0 --port 8001 --workers 4 --worker-class uvicorn.workers.UvicornWorker
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=china-car-parts-api

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/china-car-parts
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Resource limits
LimitNOFILE=65536
LimitNPROC=8192
MemoryMax=4G
CPUQuota=400%

# Health check
ExecStartPost=/bin/sleep 5
ExecStartPost=/bin/bash -c 'for i in {1..30}; do if curl -f -s http://localhost:8001/api/v1/health > /dev/null; then exit 0; fi; sleep 2; done; exit 1'

[Install]
WantedBy=multi-user.target
