# deployment/configs/china-car-parts-bot-staging.service
[Unit]
Description=China Car Parts Telegram Bot - Staging Environment
Documentation=https://github.com/your-org/china-car-parts
After=network.target postgresql.service
Wants=postgresql.service
Requires=network.target

[Service]
Type=exec
User=partsbot
Group=partsbot
WorkingDirectory=/opt/china-car-parts-staging
Environment=PATH=/opt/china-car-parts-staging/venv/bin
Environment=APP_ENV=staging
EnvironmentFile=/opt/china-car-parts-staging/.env
ExecStart=/opt/china-car-parts-staging/venv/bin/python -m app.bot.main
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=china-car-parts-bot-staging

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/china-car-parts-staging
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096
MemoryMax=1G
CPUQuota=100%

# Health check - Bot doesn't have HTTP endpoint, so we check process
ExecStartPost=/bin/sleep 10
ExecStartPost=/bin/bash -c 'for i in {1..30}; do if pgrep -f "app.bot.main" > /dev/null; then exit 0; fi; sleep 2; done; exit 1'

[Install]
WantedBy=multi-user.target
